# systemd/rotate_backups.timer
[Unit]
Description=Rotate Backups Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target