# systemd/weekly_reboot.service
[Unit]
Description=Scheduled Reboot
After=network.target

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -r now